# apps/capture-agent/ops/uk.kurdy.gtd-capture.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template — install.sh substitutes @BUN@ / @WORKDIR@ / @LABEL@ / @SCRIPT@ /
  @LOGDIR@ and writes the result to ~/Library/LaunchAgents/<label>.plist.

  KeepAlive restarts on any exit (incl. the nightly sleep/wake death that
  kept killing the Terminal-launched agent). ThrottleInterval avoids a
  crash-loop hammering the mic. RunAtLoad starts it at login.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>@LABEL@</string>

    <key>ProgramArguments</key>
    <array>
        <string>@BUN@</string>
        <string>run</string>
        <string>@SCRIPT@</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@WORKDIR@</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>@LOGDIR@/@LABEL@.out.log</string>

    <key>StandardErrorPath</key>
    <string>@LOGDIR@/@LABEL@.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:@BUNDIR@</string>
    </dict>
</dict>
</plist>
